# ansible/roles/optiora_app/templates/optiora-dashboard.service.j2
[Unit]
Description=OptiOra Dashboard
After=network.target optiora-api.service

[Service]
Type=simple
User={{ optiora_app_user }}
Group={{ optiora_app_user }}
WorkingDirectory={{ optiora_app_dir }}/dashboard
EnvironmentFile={{ optiora_app_dir }}/.env
{% if optiora_enforce_oci_runtime | default(true, true) %}
ExecStartPre=/usr/bin/bash -lc '/usr/bin/curl -fsS --max-time 3 -H "Authorization: Bearer Oracle" http://169.254.169.254/opc/v2/instance/ >/dev/null'
{% endif %}
ExecStart={{ optiora_app_dir }}/dashboard/node_modules/.bin/next start --hostname 0.0.0.0 --port {{ optiora_dashboard_port }}
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths={{ optiora_app_dir }} /var/log
LimitNOFILE={{ optiora_limits_nofile_hard }}
StandardOutput=append:/var/log/optiora-dashboard.log
StandardError=append:/var/log/optiora-dashboard.log

[Install]
WantedBy=multi-user.target
