# ansible/roles/optiora_app/templates/optiora-api.service.j2
[Unit]
Description=OptiOra API Backend
After=network.target

[Service]
Type=simple
User={{ optiora_app_user }}
Group={{ optiora_app_user }}
WorkingDirectory={{ optiora_app_dir }}
EnvironmentFile={{ optiora_app_dir }}/.env
Environment=PATH={{ optiora_app_dir }}/venv/bin
ExecStart={{ optiora_app_dir }}/venv/bin/python -m uvicorn finops_mcp.app:app --host 0.0.0.0 --port {{ optiora_api_port }}
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths={{ optiora_app_dir }} /var/log
LimitNOFILE={{ optiora_limits_nofile_hard }}
StandardOutput=append:/var/log/optiora-api.log
StandardError=append:/var/log/optiora-api.log

[Install]
WantedBy=multi-user.target
